# configs/rtc-update-pcf8523.timer
# This timer triggers daily PCF8523 RTC time updates (RPi -> RTC).
# This keeps the RTC accurate so it has good time on next boot.
#
# Installation: python simoc-sam.py setup-rtc
#
# View timer status:
#   systemctl list-timers rtc-update-pcf8523.timer
#   systemctl status rtc-update-pcf8523.timer

[Unit]
Description=Update PCF8523 RTC time from RPi time daily

[Timer]
# Run daily at 3 AM
OnCalendar=*-*-* 03:00:00
# Also run 5 minutes after boot (in case system was off at 3 AM)
OnBootSec=5min
# Run missed timers on next boot
Persistent=true
Unit=rtc-update-pcf8523.service

[Install]
WantedBy=timers.target
